# supervisord.conf
[supervisord]
nodaemon=true              
logfile=/dev/null           
pidfile=/tmp/supervisord.pid
# user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock  

[unix_http_server]
file=/tmp/supervisor.sock

[program:initialization]
command=/usr/local/bin/uv run entrypoint.py
priority=10
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:photon]
command=/usr/local/bin/start-photon.sh
process_name=photon         
autostart=false             
autorestart=false           
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
